# deploy/systemd/pycluster-retention.service
[Unit]
Description=pyCluster retention cleanup

[Service]
Type=oneshot
User=pycluster
Group=pycluster
WorkingDirectory=/home/pycluster/pyCluster
ExecStart=/usr/local/bin/pycluster-python /home/pycluster/pyCluster/scripts/cleanup_retention.py --config /home/pycluster/pyCluster/config/pycluster.toml
