# deploy/systemd/pycluster.service
[Unit]
Description=pyCluster core service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pycluster
Group=pycluster
WorkingDirectory=/home/pycluster/pyCluster
Environment=PYTHONPATH=src
ExecStart=/usr/bin/python3 -m pycluster.cli --config config/pycluster.toml serve-core
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
