# deploy/systemd/pycluster-retention.timer
[Unit]
Description=Daily pyCluster retention cleanup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h
Unit=pycluster-retention.service

[Install]
WantedBy=timers.target
